# internal/profile/add/testdata/fixtures/TestAdd/NewProfile/darwin/root/home/user/Library/LaunchAgents/local.imapnotify.OldProfile.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.imapnotify.OldProfile</string>
    <key>ProgramArguments</key>
    <array>
      <string>/home/user/.local/bin/goimapnotify</string>
      <string>-conf</string>
      <string>/home/user/.config/imapnotify/OldProfile/notify.conf</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/bin:/usr/bin:/usr/local/bin:/home/user/.local/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <true/>
    </dict>
    <key>ExitTimeOut</key>
    <integer>0</integer>
    <key>ProcessType</key>
    <string>Interactive</string>
  </dict>
</plist>
